# config/web.supervisord.conf.erb
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
pidfile=/tmp/supervisord.pid
nodaemon=true

[program:web]
command=./codecov web --port 888%(process_num)d
process_name=web_%(process_num)d
redirect_stderr=true
autorestart=true
stdout_logfile=logs/access.log
stderr_logfile=logs/error.log
numprocs=3
priority=1

[program:nginx]
command=bash -c "sleep 4 && exec bin/nginx -p . -c config/nginx.conf"
process_name=nginx
redirect_stderr=true
autorestart=true
stdout_logfile=logs/access.log
stderr_logfile=logs/error.log
priority=2
