# config/worker.supervisord.conf.erb
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
pidfile=/tmp/supervisord.pid
nodaemon=true

[program:celery]
command=./codecov worker -c 4
process_name=%(program_name)s
redirect_stderr=true
stdout_logfile=logs/access.log
stderr_logfile=logs/error.log
autorestart=true
priority=1
